# virtualhosts/etc/systemd/system/rh-php72-php-fpm.service.d/nethserver.conf
[Service]
RuntimeDirectory=rh-php72-php-fpm
RuntimeDirectoryMode=0755
