# default/etc/systemd/system/httpd.service.d/free_mutex.conf
[Service]
ExecStartPre=-/bin/bash -c "ipcs -s | awk '$3==\"apache\" { print $2 }' | xargs -r -- ipcrm sem"
